# supervisor.conf
[supervisord]
nodaemon=true

[program:sanicproxy]
command=python3 /main/main.py
stopsignal=KILL
stopasgroup=true

#[program:pythoncron]
#command=python3 /cron_job.py
#autostart=true
#autorestart=true